# cull_uid.cron
#!/bin/bash 

/usr/bin/python /opt/app/arcsight/cull_uid.py > /var/tmp/cull_uid_output